# ansible/roles/setup_dev/templates/git-check.service.j2
[Unit]
Description=Check Git Repos

[Service]
ExecStart={{ setup_dev_git_check_script_path }}

[Install]
WantedBy=multi-user.target
